# ops/systemd/chat-backup.timer
[Unit]
Description=Run Classic Chat production backup daily

[Timer]
OnCalendar=*-*-* 03:30:00
Persistent=true
RandomizedDelaySec=10m
Unit=chat-backup.service

[Install]
WantedBy=timers.target
